# extras/ImportFires/ImportFires.cron
crontab -e

0 1 * * * /usr/bin/python /ImportFires/ImportFires.py
0 4 * * * /usr/bin/python /ImportFires/ImportFires.py
0 10 * * * /usr/bin/python /ImportFires/ImportFires.py
0 13 * * * /usr/bin/python /ImportFires/ImportFires.py
0 16 * * * /usr/bin/python /ImportFires/ImportFires.py
0 19 * * * /usr/bin/python /ImportFires/ImportFires.py
0 22 * * * /usr/bin/python /ImportFires/ImportFires.py
